# flex/CONTROL/prerm
#!/bin/bash

systemctl stop platypus
systemctl disable platypus
systemctl daemon-reload

