# wsnlab/CONTROL/postinst
#!/bin/bash

rm /etc/localtime
ln -s /usr/share/zoneinfo/Europe/Paris /etc/localtime
date

systemctl daemon-reload

echo "> enable autostart service with 'systemctl enable platypus-test'"

