# platypus-test/CONTROL/prerm
#!/bin/bash

systemctl stop platypus-test
systemctl disable platypus-test
systemctl daemon-reload

