# platypus/CONTROL/prerm
#!/bin/bash

systemctl stop platypus
systemctl disable platypus

systemctl stop platypus-test
systemctl disable platypus-test

