# platypus/CONTROL/prerm
#!/bin/bash

systemctl stop platypus
systemctl disable platypus

